# ops/launchd/com.nev.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nev.daily</string>

    <!-- 由 install-daily.sh 把 REPLACE_ME 替换成 $HOME 实际值 -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPLACE_ME/nev-brief/ops/launchd/run-daily.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PROJECT_ROOT</key>
        <string>REPLACE_ME/nev-brief</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>REPLACE_ME/nev-brief</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>REPLACE_ME/nev-brief/logs/launchd-daily.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_ME/nev-brief/logs/launchd-daily.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
